# data/rauc-disk-updater.service.in
[Unit]
Description=Disk updater for Rauc
After=rauc.service local-fs.target
Requires=rauc.service

[Service]
Type=dbus
BusName=de.helbling.DiskUpdater
ExecStart=@BINDIR@/rauc-disk-updater --script @SYSCONFDIR@/rauc-disk-updater/hook.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
